# p2p-trading.service
[Unit]
Description=P2P Trading Automation System
After=network.target postgresql.service

[Service]
Type=simple
User=user
WorkingDirectory=/home/user/projects/itrader_backend
Environment="DATABASE_URL=postgresql://postgres:postgres@localhost:5432/p2p_trading"
Environment="PATH=/home/user/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/home/user/projects/itrader_backend/.venv/bin/python /home/user/projects/itrader_backend/src/main_launcher.py
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target